# build-tools/airootfs/etc/skel/.config/autostart/satpaperd.desktop
[Desktop Entry]
Exec=/usr/bin/satpaperd
Icon=application-x-shellscript
Name=satpaperd
Type=Application
X-KDE-AutostartScript=true
